# packaging/deb/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Пользователь autosshtunnels
        if ! id autosshtunnels &>/dev/null; then
            useradd -r -m -s /bin/bash autosshtunnels
        fi

        # SSH-ключ
        SSH_DIR="/home/autosshtunnels/.ssh"
        KEY_FILE="${SSH_DIR}/id_ed25519"
        if [ ! -f "$KEY_FILE" ]; then
            mkdir -p "$SSH_DIR"
            ssh-keygen -t ed25519 -f "$KEY_FILE" -N "" -C "autosshtunnels@$(hostname)"
            chown -R autosshtunnels:autosshtunnels "$SSH_DIR"
            chmod 700 "$SSH_DIR"
            chmod 600 "$KEY_FILE"
        fi

        # Директория логов
        mkdir -p /var/log/ssh-tunnel
        chown autosshtunnels:autosshtunnels /var/log/ssh-tunnel

        # После обновления: перезапуск enabled-сервисов
        for svc in $(systemctl list-unit-files --type=service --no-legend 'autossh-tunnel-*' 2>/dev/null | awk '$2=="enabled"{print $1}'); do
            systemctl start "$svc" 2>/dev/null || true
        done
        if systemctl is-enabled --quiet tunnel-watchdog.timer 2>/dev/null; then
            systemctl start tunnel-watchdog.timer 2>/dev/null || true
        fi
        ;;
esac

exit 0
